# files/_all/etc/crontabs/root
*/6 * * * * sh /etc/freifunk/ping-hack
* * * * * sh /usr/sbin/ff-stats

