# files/_all/etc/crontabs/root
*/5 * * * * sh /usr/sbin/ff-stats
0 4 * * * sh /usr/sbin/ff-fwupgrade
*/15 * * * * sh /usr/sbin/ff-rebootIfNoGw
