# files/_all/etc/hotplug.d/net/30-trafficshaping
#source:: https://github.com/freifunk-gluon/packages/blob/master/gluon/gluon-simple-tc/files/etc/hotplug.d/net/50-gluon-simple-tc              
#copyright: NeoRaider
#
#The code of Project Gluon may be distributed under the following terms, unless
#noted otherwise in indiviual files or subtrees.
#
#Copyright (c) 2013, Project Gluon
#All rights reserved.
#
#Redistribution and use in source and binary forms, with or without
#modification, are permitted provided that the following conditions are met:
#
#  1. Redistributions of source code must retain the above copyright notice,
#     this list of conditions and the following disclaimer.
#  2. Redistributions in binary form must reproduce the above copyright notice,
#     this list of conditions and the following disclaimer in the documentation
#     and/or other materials provided with the distribution.
#
#THIS SOFTWARE IS PROVIDED BY THE COPYRIGHT HOLDERS AND CONTRIBUTORS "AS IS"
#AND ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
#IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE ARE
#DISCLAIMED. IN NO EVENT SHALL THE COPYRIGHT HOLDER OR CONTRIBUTORS BE LIABLE
#FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL
#DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS OR
#SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS INTERRUPTION) HOWEVER
#CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT LIABILITY,
#OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY OUT OF THE USE
#OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF SUCH DAMAGE.

[ "$ACTION" = 'add' ] || exit 0

config_load freifunk


tc_interface() {
	local iface="$1"

	config_get ifname "$iface" ifname

	[ "$INTERFACE" = "$ifname" ] || return

	config_get_bool enabled "$iface" enabled 0

	[ "$enabled" -eq 1 ] || return

	insmod sch_ingress
	insmod sch_tbf
	insmod cls_u32
	insmod act_police

	config_get limit_out_rate "$iface" limit_out_rate
	config_get limit_in_rate "$iface" limit_in_rate

	if [ "$limit_out_rate" ]; then
		tc qdisc add dev "$INTERFACE" root tbf rate "${limit_out_rate}kbit" latency 50ms burst 2k
	fi

	if [ "$limit_in_rate" ]; then
		tc qdisc add dev "$INTERFACE" handle ffff: ingress
		tc filter add dev "$INTERFACE" parent ffff: u32 match u8 00 00 at 0 police rate "${limit_in_rate}kbit" burst 10k drop flowid :1
	fi
}

config_foreach tc_interface 'interface'
