# files/default/etc/crontabs/root
*/5  * * * * /usr/sbin/ff-stats
0    4 * * * /usr/sbin/ff-fwupgrade
*/15 * * * * /usr/sbin/ff-rebootIfNoGw
*/1  * * * * /usr/sbin/ff-offline-ssid.sh 120
